# sensorcollect.service
[Unit]
Description=Collect service of sensors
After=network-online.target

[Service]
User=INSTALL_USER
ExecStart=/usr/bin/python3 /home/INSTALL_USER/.config/station/sensorscript.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
